# docker/rdmo/rootfs/home/rdmo/conf/supervisor.conf
[program:rdmo]
command = ${HOME}/sh/runserver.sh
autorestart = true
restartpause = 1
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
